# apps/server/systemd/mysterio.service
[Unit]
Description=Mysterio mystery game server
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/mysterio/apps/server
EnvironmentFile=%h/mysterio/.env
ExecStart=/usr/bin/env node --env-file=%h/mysterio/.env dist/src/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
